# root/etc/cron.d/clamav
# Run freshclam to get new ClamAV definitions once a day
53 01 * * * root freshclam --quiet
